# dev-worker/supervisord.conf
; supervisord.conf — Waymark Dev Worker (Copilot CLI edition)
; Two programs, started in priority order:
;   10 → xvfb         (virtual display for real browser test execution)
;   20 → agent-runner (Copilot CLI watchdog loop)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
environment=HOME="/root",DISPLAY=":99",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.npm-global/bin"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── 1. Xvfb — virtual framebuffer for Playwright browser tests ───────────────
; The Copilot CLI agent itself is headless; Xvfb only exists so `npm test`
; can launch real (headed) Chromium instances for genuine E2E test coverage.
; Using DISPLAY=:99 to avoid collisions with any host display.
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
priority=10
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB

; ── 2. Agent runner — Copilot CLI watchdog loop ───────────────────────────────
; Sources /etc/agent-env.sh (written by entrypoint.sh) for AGENT_NAME,
; AGENT_COMMAND, and credential paths. Loops forever, restarting the CLI
; agent after each session ends.
; Prefer the workspace copy (live-editable without image rebuild) with
; the image copy as fallback.
[program:agent-runner]
command=bash -c 'exec "$( [ -x /workspace/dev-worker/scripts/agent-runner.sh ] && echo /workspace/dev-worker/scripts/agent-runner.sh || echo /scripts/agent-runner.sh )"'
environment=DISPLAY=":99",HOME="/root"
priority=20
autorestart=true
startsecs=0
; Stream to docker logs (stdout_logfile_maxbytes=0 disables rotation for stdout)
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
