# dev-worker/supervisord.conf
; supervisord.conf — Waymark Dev Worker container process management
; Manages 5 processes in priority order:
;   10 → xtigervnc      (virtual X11+VNC display)
;   20 → openbox        (window manager)
;   30 → novnc          (browser-based VNC client proxy)
;   50 → agent-watchdog (one-shot boot: launch VS Code, inject agent, write heartbeat, exit)
;   55 → heartbeat-loop (periodic heartbeat writer — reliable, LLM-independent)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
; Base environment for all programs
environment=HOME="/root",DISPLAY=":1",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── 1. TigerVNC — virtual X11 display with built-in VNC server ───────────────
[program:xtigervnc]
; No VNC password — ports are bound to 127.0.0.1 only (docker-compose.yml), so
; this is safe for local dev use. No passwd file complexity or DES issues.
command=Xtigervnc :1 -rfbport 5901 -SecurityTypes None -geometry 1920x1080 -depth 24 -localhost no
priority=10
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/xtigervnc.log
stderr_logfile=/var/log/supervisor/xtigervnc.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── 2. Openbox — minimal window manager ──────────────────────────────────────
; Waits for the X display to be ready before starting
[program:openbox]
command=/bin/bash -c "
    until DISPLAY=:1 xdpyinfo >/dev/null 2>&1; do sleep 1; done;
    exec openbox"
environment=DISPLAY=":1",HOME="/root"
priority=20
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB

; ── 3. noVNC — browser-accessible VNC proxy (port 6080 → 5901) ───────────────
; Waits for the VNC port to be open before starting websockify
[program:novnc]
command=/bin/bash -c "
    until nc -z localhost 5901 2>/dev/null; do sleep 1; done;
    exec websockify --web=/usr/share/novnc/ --wrap-mode=ignore 6080 localhost:5901"
priority=30
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB

; ── 4. Agent watchdog — boots the Copilot agent (one-shot) ───────────────────
; Reads AGENT_COMMAND from /etc/agent-env.sh (written by entrypoint.sh).
; On boot: launches VS Code, injects agent command, writes initial heartbeat,
; then exits. Health monitoring is handled externally by host-watchdog.sh.
[program:agent-watchdog]
command=/scripts/agent-watchdog.sh
environment=DISPLAY=":1",HOME="/root"
priority=50
autorestart=false
startsecs=0
; Log to BOTH file (for docker cp) and stdout (for docker logs)
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 5. Heartbeat loop — periodic heartbeat writer ────────────────────────────
; Reliable heartbeat mechanism that does NOT depend on the LLM agent.
; Writes heartbeat to the Workboard Heartbeat tab every 5 minutes.
; The host-side watchdog reads these to detect stale containers.
[program:heartbeat-loop]
command=/scripts/heartbeat-loop.sh
environment=HOME="/root"
priority=55
autorestart=true
startsecs=10
startretries=999
stdout_logfile=/var/log/supervisor/heartbeat-loop.log
stderr_logfile=/var/log/supervisor/heartbeat-loop.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
